# remotepc.service
[Unit]
Description=Vela
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Development/vela
Environment=START_AGENT=false
ExecStart=%h/Development/vela/.venv/bin/vela
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
